# resource/dio-monitor.service
[Unit]
Description=dio monitor

[Service]
User=root
ExecStart=dio-monitor
Restart=always
Slice=dio.slice

[Install]
WantedBy=multi-user.target